# autostart/.config/autostart/bat-cache.desktop
[Desktop Entry]
Name=bat cache
Comment=bat cache --build
Exec=bat cache --build
Type=Application
Terminal=false
Hidden=false

